# replau_ops_package/replau-daily-backup.service
[Unit]
Description=Replau LocalAPI Daily PostgreSQL Backup

[Service]
Type=oneshot
EnvironmentFile=/etc/replau-ops.env
ExecStart=/opt/replau_ops/replau_backup.sh
