# replau_ops_package/replau-stuck-monitor.timer
[Unit]
Description=Run Replau stuck monitor every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
Unit=replau-stuck-monitor.service

[Install]
WantedBy=timers.target
