# replau_ops_package/replau-conversation-retention.timer
[Unit]
Description=Run Replau WhatsApp request privacy retention daily

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
Unit=replau-conversation-retention.service

[Install]
WantedBy=timers.target
