# replau_ops_package/replau-whatsapp-watchdog.service
[Unit]
Description=Replau WhatsApp Gateway Watchdog
After=network-online.target openclaw-gateway.service
Wants=network-online.target

[Service]
Type=oneshot
User=guill
Group=guill
WorkingDirectory=/opt/replau_ops
EnvironmentFile=/etc/replau-ops.env
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/opt/replau_ops/.venv/bin/python /opt/replau_ops/replau_whatsapp_watchdog.py
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
