# replau_driver_app_package/replau-driver-app.service
[Unit]
Description=Replau Driver App
After=network.target postgrest-localapi.service
Requires=postgrest-localapi.service

[Service]
Type=simple
WorkingDirectory=/home/guill/.openclaw/workspace/replau/replau_driver_app_package
EnvironmentFile=/home/guill/.config/replau/driver.env
ExecStart=/home/guill/.local/share/replau/services/product-venv/bin/python replau_driver_app.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
