# replau_driver_app_package/replau-driver-app.service
[Unit]
Description=Replau Driver App
After=network.target postgrest-localapi.service

[Service]
Type=simple
User=guill
Group=guill
WorkingDirectory=/opt/replau_driver_app
EnvironmentFile=/etc/replau-driver-app.env
ExecStart=/opt/replau_driver_app/.venv/bin/python /opt/replau_driver_app/replau_driver_app.py
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/home/guill/.openclaw/workspace/replau_driver_documents

[Install]
WantedBy=multi-user.target
